# .config/systemd/user/tmux.service
[Unit]
Description=Tmux Terminal Multiplexer (Persistent Server)
Documentation=man:tmux(1)

[Service]
Type=forking
# Start tmux in detatched mode so systemd knows the startup sequence finished
ExecStart=/usr/bin/tmux -f %h/.tmux_outer.conf  new-session -s M -d
ExecStop=/usr/bin/tmux kill-server
Restart=on-failure

[Install]
WantedBy=default.target
